# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/data/supervisord.log
pidfile=/tmp/supervisord.pid

[program:streamlit]
command=streamlit run app.py --server.port 8501 --server.address 0.0.0.0 --server.headless true --browser.gatherUsageStats false
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/data/streamlit.log
stderr_logfile=/app/data/streamlit.err.log

[program:scheduler]
command=python scheduler.py
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/data/scheduler.log
stderr_logfile=/app/data/scheduler.err.log
